# dist/metricsd-cedfe36-amd64/DEBIAN/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create metricsd system user and group
        if ! getent group metricsd >/dev/null; then
            addgroup --system metricsd
        fi
        if ! getent passwd metricsd >/dev/null; then
            adduser --system --home /var/lib/metricsd --no-create-home \
                    --ingroup metricsd --disabled-password --shell /bin/false \
                    metricsd
        fi

        # Create runtime directory
        mkdir -p /var/lib/metricsd
        chown metricsd:metricsd /var/lib/metricsd
        chmod 750 /var/lib/metricsd

        # Set config permissions
        chown root:metricsd /etc/metricsd/config.json
        chmod 640 /etc/metricsd/config.json

        # Set certs directory permissions
        chown metricsd:metricsd /etc/metricsd/certs
        chmod 750 /etc/metricsd/certs

        # Reload systemd and enable service (but don't start)
        systemctl daemon-reload
        systemctl enable metricsd.service || true

        echo "metricsd installed successfully!"
        echo "1. Edit /etc/metricsd/config.json to configure endpoints"
        echo "2. Enable plugins by copying /usr/lib/metricsd/plugins/*.json.example to *.json"
        echo "3. Start service: sudo systemctl start metricsd"
        ;;
esac

exit 0
